# systemd/system/remote-backup.timer
[Unit]
Description=Perform remote backup four times daily

[Timer]
OnCalendar=0/6:00:00

[Install]
WantedBy=timers.target
